# bin/lakeflow.ttyd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
   <key>Label</key>
   <string>lakeflow.ttyd</string>
   <key>Program</key>
   <string>/opt/homebrew/bin/ttyd</string> 
   <key>ProgramArguments</key>
   <array>
        <string>ttyd</string>
        <string>-W</string>
        <string>/opt/homebrew/bin/tmux</string>
        <string>new</string>
        <string>-A</string>
        <string>-s</string>
        <string>Lakeflow</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>